# homebox/rootfs/etc/s6-overlay/s6-rc.d/homebox/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Homebox
# Long-running service: starts the Homebox binary with the correct environment.
# ==============================================================================

# Fix: Homebox tries to create temp files under $TMPDIR.
# /tmp is not writable in HA addon containers, so we redirect to /data/tmp.
mkdir -p /data/tmp
chmod 1777 /data/tmp
export TMPDIR="/data/tmp"

# Change to /data so Homebox resolves all relative storage paths
# inside the addon's persistent data directory.
cd /data || exit 1

# Use Homebox's official SQLite configuration variable and keep DB under
# persistent /data. Keep WAL mode enabled for SQLite durability.
export HBOX_DATABASE_SQLITE_PATH="/data/homebox.db?_pragma=busy_timeout=30000&_pragma=journal_mode=WAL&_fk=1&_time_format=sqlite"

# Show the effective SQLite database file path in startup logs.
HBOX_DATABASE_SQLITE_FILE="${HBOX_DATABASE_SQLITE_PATH%%\?*}"
mkdir -p "$(dirname "${HBOX_DATABASE_SQLITE_FILE}")"
bashio::log.info "Using SQLite database at ${HBOX_DATABASE_SQLITE_FILE}"

# Run in production mode (disables developer tooling/debug endpoints).
export HBOX_MODE="production"

# Apply user-configurable options from the HA add-on UI.
HBOX_LOG_LEVEL="$(bashio::config 'log_level')"
export HBOX_LOG_LEVEL
export HBOX_LOG_FORMAT="text"
HBOX_WEB_MAX_UPLOAD_SIZE="$(bashio::config 'max_upload_size')"
export HBOX_WEB_MAX_UPLOAD_SIZE
HBOX_OPTIONS_ALLOW_REGISTRATION="$(bashio::config 'allow_registration')"
export HBOX_OPTIONS_ALLOW_REGISTRATION

# Disable telemetry and self-update checks (not useful in a local HA context).
export HBOX_OPTIONS_ALLOW_ANALYTICS="false"
export HBOX_OPTIONS_GITHUB_RELEASE_CHECK="false"

# Trust reverse-proxy headers — always enabled because the HA Supervisor and
# the local nginx ingress proxy sit in front of Homebox.
export HBOX_OPTIONS_TRUST_PROXY="true"

bashio::log.info "Starting Homebox..."
exec /usr/bin/homebox
