# homebox/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Homebox
# Starts nginx as an ingress reverse-proxy in front of Homebox.
# ==============================================================================

declare ingress_entry
# Obtain the unique ingress path assigned by the Supervisor.
ingress_entry=$(bashio::addon.ingress_entry)
bashio::log.info "Ingress entry: ${ingress_entry}"

# Create writable runtime directory for nginx temp files and pid.
mkdir -p /tmp/nginx
chmod 700 /tmp/nginx

# Generate a self-contained nginx configuration.
# We bypass Alpine's default nginx.conf entirely so the addon does not depend
# on any distro defaults (user directive, log paths, etc.) that are not
# writable or permitted inside the restricted HA addon container.
cat > /tmp/nginx/ingress.conf <<EOF
# Run as current user (root). Without this directive nginx falls back to the
# compiled-in default user (nginx / UID 100) and calls chown / initgroups,
# which fails in the restricted HA addon container.
user root;
worker_processes 1;
error_log stderr warn;
pid /tmp/nginx/nginx.pid;

events {
    worker_connections 512;
}

http {
    include /etc/nginx/mime.types;
    default_type application/octet-stream;
    access_log off;

    client_body_temp_path /tmp/nginx/client_body;
    proxy_temp_path       /tmp/nginx/proxy;
    fastcgi_temp_path     /tmp/nginx/fastcgi;
    uwsgi_temp_path       /tmp/nginx/uwsgi;
    scgi_temp_path        /tmp/nginx/scgi;

    sendfile on;
    tcp_nopush on;
    keepalive_timeout 65;

    map \$http_upgrade \$connection_upgrade {
        default upgrade;
        ''      close;
    }

    server {
        listen 8099 default_server;

        # Only the HA Supervisor may connect to the ingress port.
        # 172.30.32.2 is the Supervisor's fixed IP on the hassio bridge network.
        allow   172.30.32.2;
        deny    all;

        server_name _;

        proxy_http_version          1.1;
        proxy_ignore_client_abort   on;
        proxy_buffering             off;
        proxy_read_timeout          3600s;
        proxy_send_timeout          3600s;

        location / {
            proxy_pass              http://127.0.0.1:7745;

            proxy_set_header Host              \$http_host;
            proxy_set_header X-Real-IP         \$remote_addr;
            proxy_set_header X-Forwarded-For   \$proxy_add_x_forwarded_for;
            proxy_set_header X-Forwarded-Proto \$scheme;
            proxy_set_header Upgrade           \$http_upgrade;
            proxy_set_header Connection        \$connection_upgrade;

            # Strip headers that prevent the page from rendering inside the
            # HA ingress iframe.
            proxy_hide_header X-Frame-Options;
            proxy_hide_header Content-Security-Policy;
            proxy_hide_header Content-Origin-Embedder-Policy;
            proxy_hide_header Content-Origin-Opener-Policy;
            proxy_hide_header Content-Origin-Resource-Policy;

            # Fix cookie Domain: the HA Supervisor replaces the Host header
            # with the addon's internal IP, so Homebox sets cookies with
            # Domain=<internal-ip> which the browser rejects (origin
            # mismatch).  Stripping the Domain attribute lets the browser
            # scope cookies to the actual request origin.
            proxy_cookie_domain ~.+ "";

            # Allow file uploads up to the configured maximum size.
            client_max_body_size 100m;

            # Rewrite root-relative URLs that Homebox embeds in HTML / JS / CSS
            # so that they work under the HA ingress sub-path.
            sub_filter_once off;
            sub_filter_types text/css text/javascript application/javascript application/json application/manifest+json;
            sub_filter 'baseURL:"/"'  'baseURL:"${ingress_entry}/"';
            sub_filter 'src="/'       'src="${ingress_entry}/';
            sub_filter 'href="/'      'href="${ingress_entry}/';
            sub_filter 'src:"/'       'src:"${ingress_entry}/';
            sub_filter 'href:"/'      'href:"${ingress_entry}/';
            sub_filter "src:'/"       "src:'${ingress_entry}/";
            sub_filter "href:'/"      "href:'${ingress_entry}/";
            sub_filter 'url(/'        'url(${ingress_entry}/';
            sub_filter 'url("/'       'url("${ingress_entry}/';
            sub_filter "url('/"       "url('${ingress_entry}/";
            sub_filter '}/api/'       '}${ingress_entry}/api/';
            sub_filter '"/api/'       '"${ingress_entry}/api/';
            sub_filter "'/api/"       "'${ingress_entry}/api/";
            sub_filter ':"/no-image.jpg"' ':"${ingress_entry}/no-image.jpg"';
            sub_filter ":'/no-image.jpg'" ":'${ingress_entry}/no-image.jpg'";

            proxy_redirect  / ${ingress_entry}/;
        }
    }
}
EOF

bashio::log.info "Starting nginx ingress proxy..."
exec nginx -c /tmp/nginx/ingress.conf -e stderr -g 'daemon off; master_process off;'
